# deploy/systemd/hermes-watchdog.service
[Unit]
Description=Hermes - Watchdog (heartbeat, config drift, eventi)

[Service]
Type=oneshot
WorkingDirectory=%h/hermes_hft
EnvironmentFile=-%h/hermes_hft/.env
ExecStart=%h/hermes_hft/venv/bin/python watchdog.py
StandardOutput=append:%h/hermes_hft/logs/watchdog.log
StandardError=append:%h/hermes_hft/logs/watchdog.log
